# vision-usb-gateway/systemd/citostore-usb-mount@.service
[Unit]
Description=Mount USB export drive %I
# udev pulls this in via SYSTEMD_WANTS when a USB filesystem appears.
BindsTo=dev-%i.device
After=dev-%i.device

[Service]
Type=oneshot
RemainAfterExit=no
Environment=GATEWAY_HOME=/opt/CitoStore/vision-usb-gateway
EnvironmentFile=-/etc/vision-gw.env
ExecStart=/bin/bash -c '"${GATEWAY_HOME:-/opt/CitoStore/vision-usb-gateway}"/scripts/usb-export-mount.sh add %I'
