# vision-usb-gateway/systemd/citostore-mdns.service
[Unit]
Description=CitoStore mDNS name mode (network vs direct 1-1)
# Evaluate once addressing has settled so a slow DHCP can't be seen as "direct".
After=network-online.target avahi-daemon.service vision-gw-config.service
Wants=network-online.target
Requires=avahi-daemon.service

[Service]
Type=oneshot
RemainAfterExit=yes
Environment=GATEWAY_HOME=/opt/CitoStore/vision-usb-gateway
ExecStart=/bin/bash ${GATEWAY_HOME}/scripts/mdns-apply-mode.sh

[Install]
WantedBy=multi-user.target
