# vision-usb-gateway/systemd/nas-sync.timer
[Unit]
Description=NAS Sync Timer

[Timer]
OnUnitActiveSec=10min
Persistent=true

[Install]
WantedBy=timers.target