# vision-usb-gateway/systemd/usb-gadget.service
[Unit]
Description=USB Gadget for Vision Host
After=local-fs.target

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=/etc/vision-gw.env
ExecStart=/bin/bash ${GATEWAY_HOME}/scripts/usb-gadget.sh start
ExecStop=/bin/bash ${GATEWAY_HOME}/scripts/usb-gadget.sh stop
ProtectSystem=strict
PrivateTmp=true
ReadWritePaths=/run /sys/kernel/config
NoNewPrivileges=true
CapabilityBoundingSet=CAP_SYS_ADMIN CAP_SYS_RAWIO

[Install]
WantedBy=multi-user.target