# vision-usb-gateway/systemd/usb-gadget.service
[Unit]
Description=USB Gadget for Vision Host
# Order after the config restore so /etc/vision-gw.conf (and USB_LVS) is fully
# written before we start: with only After=local-fs.target the gadget raced
# vision-gw-config's rewrite of that file and intermittently read it truncated,
# crashing on `USB_LVS: unbound variable`.
After=local-fs.target vision-gw-config.service
Wants=vision-gw-config.service

[Service]
Type=oneshot
TimeoutStartSec=3min
RemainAfterExit=yes
Environment=GATEWAY_HOME=/opt/CitoStore/vision-usb-gateway
EnvironmentFile=-/etc/vision-gw.env
ExecStart=/bin/bash ${GATEWAY_HOME}/scripts/usb-gadget.sh start
ExecStop=/bin/bash ${GATEWAY_HOME}/scripts/usb-gadget.sh stop
ProtectSystem=strict
PrivateTmp=true
ReadWritePaths=/run /sys/kernel/config /srv/vision_mirror
NoNewPrivileges=true
CapabilityBoundingSet=CAP_SYS_ADMIN CAP_SYS_RAWIO

[Install]
WantedBy=multi-user.target
