# vision-usb-gateway/systemd/vision-gw-config.service
[Unit]
Description=Apply vision-gw.conf to systemd overrides
After=local-fs.target

[Service]
Type=oneshot
Environment=GATEWAY_HOME=/opt/vision-usb-gateway
EnvironmentFile=-/etc/vision-gw.env
ExecStart=/bin/bash ${GATEWAY_HOME}/scripts/update-config.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
