# vision-usb-gateway/systemd/vision-provision.service
[Unit]
Description=CitoStore provision from config bundle (DESTRUCTIVE: wipes NVMe)
After=local-fs.target
# Triggered on demand by the WebUI; never enabled at boot.

[Service]
Type=oneshot
EnvironmentFile=-/etc/vision-gw.env
ExecStart=/bin/bash -c '"${GATEWAY_HOME:-/opt/vision-usb-gateway}"/scripts/provision-from-bundle.sh /run/vision-provision/bundle.citostore --provision --confirm'
TimeoutStartSec=1800
