# vision-usb-gateway/systemd/vision-nvme-health.timer
[Unit]
Description=Vision NVMe SMART Cache Timer

[Timer]
OnBootSec=2min
OnUnitActiveSec=10min

[Install]
WantedBy=timers.target
