# vision-usb-gateway/systemd/vision-sync.timer
[Unit]
Description=Vision Snapshot Sync Timer

[Timer]
OnBootSec=2min
OnUnitActiveSec=2min
Persistent=true

[Install]
WantedBy=timers.target
