# vision-usb-gateway/systemd/vision-monitor.service
[Unit]
Description=Vision Monitor
After=local-fs.target

[Service]
Type=oneshot
EnvironmentFile=/etc/vision-gw.env
ExecStart=/bin/bash ${GATEWAY_HOME}/scripts/vision-monitor.sh
ProtectSystem=strict
PrivateTmp=true
ReadWritePaths=/run
NoNewPrivileges=true
CapabilityBoundingSet=CAP_SYS_ADMIN

[Install]
WantedBy=multi-user.target