# vision-usb-gateway/systemd/vision-monitor.timer
[Unit]
Description=Vision Monitor Timer

[Timer]
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target