# vision-usb-gateway/systemd/vision-factory-reset.service
[Unit]
Description=Vision Factory Reset (arm: mark for NVMe wipe on next boot, then reboot)
After=local-fs.target

[Service]
Type=oneshot
Environment=GATEWAY_HOME=/opt/CitoStore/vision-usb-gateway
EnvironmentFile=-/etc/vision-gw.env
# Only arms and reboots — the wipe itself runs at early boot, when the NVMe is not
# yet held. Runs as full root, outside the WebUI's sandbox; reboots on success.
ExecStart=/bin/bash ${GATEWAY_HOME}/scripts/factory-reset.sh --arm
