# vision-usb-gateway/systemd/vision-rotator.timer
[Unit]
Description=Vision Rotator Timer

[Timer]
OnUnitActiveSec=1min
Persistent=true

[Install]
WantedBy=timers.target