# vision-usb-gateway/systemd/vision-journal-persist.timer
[Unit]
Description=Persist Journal Tail Timer

[Timer]
OnBootSec=3min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
