# vision-usb-gateway/systemd/vsftpd.service.d/vision-ingest.conf
[Unit]
# vsftpd binds listen_address on the ingest interface (eth1). The packaged unit
# only orders itself After=network.target, which does not wait for an address,
# so at boot vsftpd exits 2/INVALIDARGUMENT before NetworkManager has configured
# eth1. The address is also absent for as long as the AOI's cable has no
# carrier, so retry indefinitely instead of giving up at the default start limit.
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Restart=on-failure
RestartSec=5s
