# vision-usb-gateway/systemd/avahi-daemon.service.d/vision-restart.conf
# avahi is what makes AOI1.local reachable in the field; the packaged unit has
# no Restart=, so a crash would silently take name resolution down for good.
[Service]
Restart=on-failure
RestartSec=5
