# vision-usb-gateway/systemd/NetworkManager-wait-online.service.d/timeout.conf
# Cap NetworkManager-wait-online so a not-yet-ready interface (typically eth1,
# the AOI point-to-point link with no carrier at boot) cannot stall boot for the
# packaged ~30-90s default. 15s is enough for eth0 (SMB) to come up; anything
# slower activates later on its own via autoconnect.
[Service]
ExecStart=
ExecStart=/usr/bin/nm-online -s -q --timeout=15
