# vision-usb-gateway/systemd/nmbd.service.d/network.conf
# nmbd needs a non-loopback IPv4 interface or it spins "waiting for interface"
# and times out (~90s), failing the boot and (via 50_configure_samba's restart)
# the whole apply-shadow-config. A unit may boot before its network cable is
# plugged in, so: start after network-online, cap the start timeout, and retry
# so nmbd comes up on its own once an interface appears.
[Unit]
After=network-online.target
Wants=network-online.target

[Service]
TimeoutStartSec=25
Restart=on-failure
RestartSec=10
