# vision-usb-gateway/systemd/nmbd.service.d/vision-gw.conf
[Unit]
After=vision-gw-config.service
Requires=vision-gw-config.service
